# tmproot/etc/profile.d/proxy.sh
# this file was generated with and might get overwritten by setup-proxy

export http_proxy=http://proxy.nor.wtbts.net:8080
export https_proxy=http://proxy.nor.wtbts.net:8080
export ftp_proxy=http://proxy.nor.wtbts.net:8080
